# launchagents/com.assistant.assistant-comms.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.assistant.assistant-comms</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/mukuls/dev/assistant/bin/comms-listen.py</string>
    </array>
    <!-- Event-driven daemon: a single long-running process (Telegram long-poll
         + ledger watch + heartbeat timer). KeepAlive relaunches it if it
         crashes or exits. No StartInterval — it does not tick, it listens. -->
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/mukuls/.architect/orchestrator-logs/assistant-comms.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/mukuls/.architect/orchestrator-logs/assistant-comms.launchd.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mukuls</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/cmux.app/Contents/Resources/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mukuls</string>
</dict>
</plist>
