# launchagents/com.mukuls.assistant-todo-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mukuls.assistant-todo-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/mukuls/dev/assistant/bin/todo-server.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/mukuls</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mukuls/.claude/cache/todo-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mukuls/.claude/cache/todo-server.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/mukuls</string>
</dict>
</plist>
