# launchagents/com.mukul.assistant-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.mukul.assistant-daemon — the single-process Assistant daemon.

  Collapses the com.assistant.assistant-pulse LaunchAgent (a 5-min timer) into
  ONE process: `python -m assistant`, run via bin/assistant-daemon.py. The
  daemon owns the pulse loop, the daemon heartbeat, and tool dispatch — each in
  its own thread of one process. No StartInterval: the daemon does not tick, it
  runs; KeepAlive relaunches it if it ever exits.

  THIS PLIST IS NOT LOADED BY THE INSTALLER. The migration is additive — the
  existing pulse agent keeps running until you switch over by hand.

  To activate (hand off pulse to the single daemon):
    launchctl unload ~/Library/LaunchAgents/com.assistant.assistant-pulse.plist
    launchctl load   ~/Library/LaunchAgents/com.mukul.assistant-daemon.plist

  To roll back:
    launchctl unload ~/Library/LaunchAgents/com.mukul.assistant-daemon.plist
    launchctl load   ~/Library/LaunchAgents/com.assistant.assistant-pulse.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mukul.assistant-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/mukuls/dev/assistant/bin/assistant-daemon.py</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/mukuls/.architect/orchestrator-logs/assistant-daemon.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/mukuls/.architect/orchestrator-logs/assistant-daemon.launchd.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/mukuls</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Applications/cmux.app/Contents/Resources/bin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/mukuls</string>
</dict>
</plist>
